# supervisord.conf
[supervisord]
nodaemon=true

[program:httpd]
command=service httpd start

[program:mysqld]
command=service mysqld start

[program:mongod]
command=service mongod start
